# com.jake.manager-mccode.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jake.manager-mccode</string>
    <key>ProgramArguments</key>
    <array>
        <string>/path/to/poetry</string>
        <string>run</string>
        <string>python</string>
        <string>run.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/path/to/manager_mccode</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/manager-mccode.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/manager-mccode.error.log</string>
</dict>
</plist> 